# examples/systemd/mvp-vpn-lite-client.service
[Unit]
Description=mvp-vpn-lite QUIC TUN client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/mvp-vpn-lite/client.env
ExecStartPre=/usr/local/share/mvp-vpn-lite/scripts/setup-client.sh
ExecStart=/usr/local/bin/mvp-vpn-lite-client
ExecStopPost=/usr/local/share/mvp-vpn-lite/scripts/cleanup-client.sh
Restart=on-failure
RestartSec=2s
KillSignal=SIGINT
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
